# linux/cleanup.service.in
[Unit]
Description=fwupdate housekeeping

[Service]
PrivateTmp=true
Type=oneshot
ExecStart=@@LIBEXECDIR@@fwupdate/cleanup
